# ops/com.ryan.health-dashboard.sync-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ryan.health-dashboard.sync-daily</string>

  <key>WorkingDirectory</key>
  <string>/Users/ryan-desktop/.openclaw/workspace/health-dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>npm run sync scheduled</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/ryan-desktop/.openclaw/workspace/health-dashboard/logs/health-dashboard-sync.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ryan-desktop/.openclaw/workspace/health-dashboard/logs/health-dashboard-sync.stderr.log</string>
</dict>
</plist>
